# init/kodi-logger.timer
[Unit]
Description=Hourly kodi logger

[Timer]
OnActiveSec=0m
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
