# init/kodi-logger.service
[Unit]
Description=Kodi logger

[Service]
Type=simple
ExecStart=/usr/bin/kodi-logger

[Install]
WantedBy=default.target
